# src/main/resources/jpostgrey.service
[Unit]
Description=jpostgrey
After=syslog.target

[Service]
User=root
ExecStart=/opt/jpostgrey/jpostgrey

[Install]
WantedBy=multi-user.target